# minecraft_server.service
[Unit]
Description=Minecraft virtualized server
target=

[Service]


[Install]
